# background/schedule_manifest.yaml
# Schedule / service manifest — the committed declaration of OS-LEVEL scheduling and services.
# OPS1 sub-step 3 (G-L2 applied to scheduling), docs/design/OPERATIONAL_LAYER_DESIGN.md §5.
#
# The mandate's canonical violation was "tonight's invisible cron that resurrected a broken
# stack" — an OS crontab entry nobody reading the repo could see. This manifest makes OS-level
# scheduling/services DECLARED in git, and background/schedule_reconciler.py reconciles the
# machine's ACTUAL crontab + systemd units against it: ANY undeclared crontab line or SE systemd
# unit is flagged as drift. The invisible-cron class is now structurally detectable, forever.
#
# EVERY non-trivial entry carries its `reason` (reason-with-state rule; IaC).
version: 1

# Declared user-crontab entries. EMPTY BY DESIGN — nothing is scheduled via hand-run `crontab -e`.
# The stack is started explicitly (start_worker.sh); periodic behaviour lives IN the daemons
# (committed code), not in cron. Any crontab line not listed here is UNDECLARED drift — the
# exact invisible-cron the blackout came from.
cron: []

# Declared systemd (--user) units owned by synthetic-enterprise. install_schedule.sh installs
# each from its committed unit_file (reconstruct-from-repo), so the machine's service config is
# reproducible from the repo alone.
systemd_units:
  - name: file-api.service
    unit_file: background/file-api.service
    enabled: true            # boot-start
    active: true             # must be running
    health: "http://localhost:8765/health"
    reason: >-
      file-api is DIRECTOR-ACCESS INFRASTRUCTURE, not an autonomous-company daemon: it serves
      the director independent of whether the stack is up, held, or torn down (2026-07-17: stack
      fully down, file-api still his phone-only bridge). It must be the most resilient thing on
      the box — boot-start + Restart=on-failure — because its death strands him exactly when
      he is phone-only. The split model (systemd for this, tmux for the stack) is principled,
      not messy. Owning it here is why it is NOT in process_manifest.yaml.

  - name: boot-announce.service
    unit_file: background/boot-announce.service
    enabled: true            # boot-start: fire once on every boot
    active: false            # Type=oneshot — runs, reports, EXITS; inactive(dead) is correct, not UNIT_DOWN
    reason: >-
      OPS1 sub-step 4, G-R3 (reconcile before act). A oneshot that on every boot reconciles the
      declared set (process_manifest + this manifest) against actual and sends ONE boot NTFY
      BEFORE any autonomous work advances — so a boot state is a stated single transition, never
      the silent "cron resurrects whatever was there" mask the blackout came from. Report-only:
      it starts/stops/reaps nothing (G-R4). enabled (boot-start) but active:false because a
      oneshot exits after running — flagging it UNIT_DOWN would be wrong.

  - name: reconcile-watch.service
    unit_file: background/reconcile-watch.service
    enabled: false           # timer-triggered, not boot-enabled itself (static)
    active: false            # Type=oneshot — runs one reconcile pass and EXITS
    reason: >-
      OPS1 sub-step 4: the drift CONTROL made LIVE. The oneshot the timer fires — reconciles
      declared-vs-actual and pages ONLY on a drift transition (R5), report-only (G-R3). Not
      enabled/active itself (the .timer drives it); flagging it down would be wrong.

  - name: reconcile-watch.timer
    unit_file: background/reconcile-watch.timer
    enabled: true            # boot-start (WantedBy=timers.target)
    active: true             # must be ARMED — a drift control that isn't scheduled is fail-silent (R15)
    reason: >-
      OPS1 sub-step 4: the SCHEDULE that makes the reconcile a real control instead of theatre.
      boot_announce runs the reconcile once at boot; between boots drift was UNWATCHED — the exact
      gap that let a live worker-seat declared `held` produce no HELD_VIOLATED (2026-07-17). This
      timer fires the reconcile every 5 min so drift is loud within minutes, not at the next boot.
      active:true because an un-armed timer is a control that never runs.

  # ── Scheduled bounded invocations — the continuity transport (SCHEDULED_BOUNDED_INVOCATIONS_DESIGN.md)
  # DARK until the director-run cutover (background/cutover_to_scheduled.sh --apply). Declared here so
  # reconcile SEES them; enabled/active:false = installed-not-armed, which matches reality pre-cutover.
  # At cutover these flip to enabled/active (committed as the IaC record) and worker-seat-manager is
  # retired. This replaces the persistent stop-hook-rearm seat, which failed three ways in ~30h
  # (input-blocking / death-on-API-error / token-burn at rest).
  - name: worker-tick.service
    unit_file: background/worker-tick.service
    enabled: false           # Type=oneshot fired by the timer/path — never boot-enabled itself
    active: false            # oneshot: inactive(dead) between ticks is correct, not UNIT_DOWN
    reason: >-
      The bounded-invocation check: cheap Python (origin-sync + find_work) that spawns ONE
      `claude -p` worker ONLY when there is work, and exits at ~0 tokens at rest. DARK until cutover
      (enabled/active:false = installed, not armed). Not enabled/active itself — the .timer and .path
      drive it; flagging a oneshot down would be wrong.

  - name: worker-tick.timer
    unit_file: background/worker-tick.timer
    enabled: true            # LIVE (cutover 2026-07-20): boot-start, WantedBy=timers.target
    active: true             # LIVE: armed — an un-armed continuity timer is fail-silent (R15)
    reason: >-
      The 60s wake backstop + origin-staged-doc path. LIVE since the 2026-07-20 cutover. This is the
      continuity heartbeat's replacement: level-triggered (did the timer fire?) instead of the
      edge-triggered stop-hook re-arm that could not stay alive (3 failures in ~30h: input-block,
      death-on-API-error, token-burn). Independent of worker-tick.path (R15: either trigger alone
      still wakes).

  - name: worker-tick.path
    unit_file: background/worker-tick.path
    enabled: true            # LIVE (cutover 2026-07-20): boot-start, WantedBy=paths.target
    active: true             # LIVE
    reason: >-
      The fast local-staged-doc wake (~seconds vs the 60s timer). LIVE since the 2026-07-20 cutover.
      A staged doc is the primary continuity case (both 2026-07-19 failures hit it). Independent of
      worker-tick.timer (R15: kill the timer and this still wakes on a staged doc).
